# deploy/tg-proxy.service
[Unit]
Description=RussVodokanal Telegram notify (tg-proxy.js)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/var/www/rusvodokanal
ExecStart=/usr/bin/node /var/www/rusvodokanal/tg-proxy.js
Restart=always
RestartSec=5
# Ограничить пользователя (нужны права на чтение node_modules):
# User=www-data
# Group=www-data
NoNewPrivileges=true

# Секреты лучше вынести сюда:
# EnvironmentFile=-/etc/default/tg-proxy

[Install]
WantedBy=multi-user.target
